# examples/dynamic-dns.service
[Unit]
Description=Dynamic DNS service

[Service]
Type=oneshot
ExecStart=/opt/dynamic-dns/dynamic-dns
