# examples/dynamic-dns.timer
[Unit]
Description=Update DNS every 15 minutes

[Timer]
OnBootSec=60s
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target
